# integrations/load-forecasting/deploy/aether-load-forecasting-processor.service
[Unit]
Description=Aether request-driven load forecasting processor
After=network.target

[Service]
Type=simple
EnvironmentFile=/etc/aether/load-forecasting-processor.env
WorkingDirectory=/opt/load-forecasting
ExecStart=/opt/load-forecasting/.venv/bin/uvicorn app:app --host 127.0.0.1 --port 8989
Restart=on-failure
RestartSec=5
DynamicUser=yes
NoNewPrivileges=yes
PrivateDevices=yes
PrivateTmp=yes
ProtectControlGroups=yes
ProtectHome=yes
ProtectKernelModules=yes
ProtectKernelTunables=yes
ProtectSystem=strict
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
RestrictNamespaces=yes
RestrictRealtime=yes
SystemCallArchitectures=native
UMask=0077

[Install]
WantedBy=multi-user.target
